# roles/install-cockpit/files/etc/systemd/system/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=192.168.1.2:7070
FreeBind=yes
